# systemd/alphacouncil-worker.service
[Unit]
Description=AlphaCouncil Worker (Mantle Indexer + AI Pipeline)
Documentation=https://github.com/your-org/alphacouncil
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=alphacouncil
Group=alphacouncil
WorkingDirectory=/opt/alphacouncil/apps/worker
EnvironmentFile=/opt/alphacouncil/.env
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alphacouncil-worker

[Install]
WantedBy=multi-user.target
